# crontab
*/1 * * * * /bin/bash /root/cloudflare-ddns-updater/cloudflare.sh
